# src/etc/local.d/2_user_setup.start
#!/bin/sh
echo "setxkbmap de_CH &" > /home/develop/.xinitrc
echo "exec xfce4-session" >> /home/develop/.xinitrc
chown develop:develop /home/develop/.xinitrc
