# src/etc/local.d/3_fix_permissions.start
#!/bin/sh
chown -R develop:develop /home/develop
chown -R develop:develop /tmp
